# run/systemd/generator/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=local-fs.target

[Mount]
Where=/
What=/dev/disk/by-uuid/2b02dda1-55f1-4fbd-ba68-cfb90ae97790
Type=ext4
